# nekobrowser.service
[Unit]
Description=NekoBrowser Auto-Kiosk Service
After=network.target

[Service]
Type=simple
ExecStart=/usr/bin/python3 /home/username/NekoBrowser/neko_enhanced.py
WorkingDirectory=/home/username/NekoBrowser
Restart=always
RestartSec=5
User=username
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target
